# systemd/musick-cache.service
[Unit]
Description=Musick cache updater

[Service]
Type=oneshot
ExecStart=@@INSTALL_DIR@@/media-cache.sh
